# debian/postinst
#!/bin/sh

set -e
umask 0022

. /usr/share/debconf/confmodule

update-rc.d lightldapd defaults 90 > /dev/null
invoke-rc.d lightldapd start

db_stop
