# debian/prerm
#!/bin/sh -e

if [ -x /etc/init.d/lightldapd ]; then
	invoke-rc.d lightldapd stop
fi
